# debian/DEBIAN/prerm
#!/bin/sh
systemctl stop tupass
systemctl disable tupass